# templates/sync-replica-incoming-watch-messages.path
[Unit]
Description=Sync replica incoming watch messages

[Path]
PathChanged=/opt/sync/replica-incoming/transfer-messages.event

[Install]
WantedBy=multi-user.target
