# templates/sync-replica-send-watch.path
[Unit]
Description=Sync replica send watch

[Path]
PathChanged=/run/sync-replica.event

[Install]
WantedBy=multi-user.target
