# templates/sync-replica-incoming-watch.path
[Unit]
Description=Sync replica incoming watch

[Path]
PathChanged=/opt/sync/replica-incoming/transfer.event

[Install]
WantedBy=multi-user.target
